# supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:wes]
command=nice -10 python3 /usr/src/app/wes/main.py -rcsv /usr/src/app/projects.csv
user=wes
environment=HOME="/home/wes",USER="wes"
numprocs=1
autostart=true
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
command=/usr/local/bin/gunicorn -c /usr/src/app/gunicorn.conf -u wes web:app
directory=/usr/src/app/wes
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=crond -f
autostart=true
autorestart=false
